# docker/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:auth-service]
command=dotnet /app/auth/AuthService.Api.dll
directory=/app/auth
autostart=true
autorestart=true
priority=10
environment=ASPNETCORE_URLS="http://127.0.0.1:8081",ASPNETCORE_ENVIRONMENT="Production",ConnectionStrings__AuthDb="Server=127.0.0.1;Database=TokaAuth;User Id=sa;Password=Toka@123456;TrustServerCertificate=True;",Messaging__UseRabbitMq="true",RabbitMq__Host="127.0.0.1",Jwt__Issuer="toka-auth",Jwt__Audience="toka-services",Jwt__SecretKey="TokaSuperSecretKeyForJwtTokens2024!MustBeLongEnough",Jwt__ExpirationMinutes="60"
stdout_logfile=/var/log/supervisor/auth.log
stderr_logfile=/var/log/supervisor/auth.err

[program:user-service]
command=dotnet /app/user/UserService.Api.dll
directory=/app/user
autostart=true
autorestart=true
priority=10
environment=ASPNETCORE_URLS="http://127.0.0.1:8082",ASPNETCORE_ENVIRONMENT="Production",ConnectionStrings__UserDb="Server=127.0.0.1;Database=TokaUsers;User Id=sa;Password=Toka@123456;TrustServerCertificate=True;",ConnectionStrings__Redis="127.0.0.1:6379",Caching__UseRedis="true",Messaging__UseRabbitMq="true",RabbitMq__Host="127.0.0.1",Jwt__Issuer="toka-auth",Jwt__Audience="toka-services",Jwt__SecretKey="TokaSuperSecretKeyForJwtTokens2024!MustBeLongEnough",Jwt__ExpirationMinutes="60"
stdout_logfile=/var/log/supervisor/user.log
stderr_logfile=/var/log/supervisor/user.err

[program:role-service]
command=dotnet /app/role/RoleService.Api.dll
directory=/app/role
autostart=true
autorestart=true
priority=10
environment=ASPNETCORE_URLS="http://127.0.0.1:8083",ASPNETCORE_ENVIRONMENT="Production",ConnectionStrings__RoleDb="Server=127.0.0.1;Database=TokaRoles;User Id=sa;Password=Toka@123456;TrustServerCertificate=True;",Messaging__UseRabbitMq="true",RabbitMq__Host="127.0.0.1",Jwt__Issuer="toka-auth",Jwt__Audience="toka-services",Jwt__SecretKey="TokaSuperSecretKeyForJwtTokens2024!MustBeLongEnough",Jwt__ExpirationMinutes="60"
stdout_logfile=/var/log/supervisor/role.log
stderr_logfile=/var/log/supervisor/role.err

[program:audit-service]
command=dotnet /app/audit/AuditService.Api.dll
directory=/app/audit
autostart=true
autorestart=true
priority=10
environment=ASPNETCORE_URLS="http://127.0.0.1:8084",ASPNETCORE_ENVIRONMENT="Production",MongoDb__ConnectionString="mongodb://127.0.0.1:27017",MongoDb__UseMongo="true",MongoDb__DatabaseName="TokaAudit",Messaging__UseRabbitMq="true",RabbitMq__Host="127.0.0.1",Jwt__Issuer="toka-auth",Jwt__Audience="toka-services",Jwt__SecretKey="TokaSuperSecretKeyForJwtTokens2024!MustBeLongEnough",Jwt__ExpirationMinutes="60"
stdout_logfile=/var/log/supervisor/audit.log
stderr_logfile=/var/log/supervisor/audit.err

[program:ai-agent-service]
command=dotnet /app/agent/AiAgentService.Api.dll
directory=/app/agent
autostart=true
autorestart=true
priority=10
environment=ASPNETCORE_URLS="http://127.0.0.1:8085",ASPNETCORE_ENVIRONMENT="Production",Qdrant__Enabled="false",Services__UserService="http://127.0.0.1:8082/",Llm__Provider="%(ENV_Llm__Provider)s",Ollama__Enabled="%(ENV_Ollama__Enabled)s",Ollama__BaseUrl="%(ENV_Ollama__BaseUrl)s",Ollama__ChatModel="%(ENV_Ollama__ChatModel)s",Ollama__EmbeddingModel="%(ENV_Ollama__EmbeddingModel)s",Ollama__EmbeddingDimensions="%(ENV_Ollama__EmbeddingDimensions)s",OpenAi__ApiKey="%(ENV_OPENAI_API_KEY)s"
stdout_logfile=/var/log/supervisor/agent.log
stderr_logfile=/var/log/supervisor/agent.err

[program:gateway]
command=dotnet /app/gateway/Gateway.Api.dll
directory=/app/gateway
autostart=true
autorestart=true
priority=20
startsecs=3
environment=ASPNETCORE_URLS="http://127.0.0.1:8080",ASPNETCORE_ENVIRONMENT="Production"
stdout_logfile=/var/log/supervisor/gateway.log
stderr_logfile=/var/log/supervisor/gateway.err

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
